# config/systemd/user/inbox.timer
[Unit]
Description=Sync ~/dev/inbox/ with GitHub every minute

[Timer]
OnCalendar=*:0/1
Persistent=true

[Install]
WantedBy=timers.target
