# config/systemd/user/ssh-agent.service
[Unit]
Description=ssh-agent
Requires=ssh-agent.socket

[Service]
ExecStart=/usr/bin/ssh-agent -d -a %t/ssh_auth_sock
ExecStop=/usr/bin/ssh-agent -k
